# dotfiles/.config/systemd/user/battery-ac-check.service
[Unit]
Description=Battery AC Status Check

[Service]
Type=oneshot
ExecStart=%h/.scripts/battery-ac-check.sh
Environment=DISPLAY=:0
