# dotfiles/.config/systemd/user/mango-reload.service
[Unit]
Description=Reload Mango Window Manager

[Service]
Type=oneshot
ExecStart=mmsg -d reload_config
