# dotfiles/.config/systemd/user/mango-reload.path
[Unit]
Description=Watch Mango config file for changes

[Path]
PathModified=%h/.config/mango/

[Install]
WantedBy=default.target
